# laradock/workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock /var/www/laradock/check2.sh >> /tmp/check.log 2>&1
* * * * * laradock echo "cron start" > /tmp/log.log 2>&1
